# deploy/agent-passport.service.example
[Unit]
Description=agent-passport
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agentpassport
Group=agentpassport
WorkingDirectory=/opt/agent-passport/current
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/agent-passport/agent-passport.env
ExecStart=/usr/bin/env node src/server.js
Restart=always
RestartSec=5
StateDirectory=agent-passport
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/agent-passport
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
